# templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
[Service]
Type=simple
ExecStart={{ ne_dir }}/node_exporter
[Install]
WantedBy=multi-user.target